# deploy/alas-launcher.service
[Unit]
Description=ALAS Launch Button Watcher (physical boot/standby button)
Documentation=file://%h/AI-Glasses/src/main/boot_launcher.py
After=multi-user.target sound.target
Wants=multi-user.target

[Service]
Type=simple

# ─────────────────────────────────────────────────────────────────────────
# User        : the account that owns the repo / venv (NOT root if avoidable).
# WorkingDir  : the repo's  src/  directory (so `-m main.boot_launcher` resolves).
# ExecStart   : the venv's python running the launcher module.
User=alas
WorkingDirectory=/home/alas/ALAS_PROJECT/AI-Glasses/src
ExecStart=/home/alas/ALAS_PROJECT/AI-Glasses/alas_env/bin/python -m main.boot_launcher
# ─────────────────────────────────────────────────────────────────────────

# GPIO access needs the 'gpio' group; audio output (TTS) needs 'audio'.
SupplementaryGroups=gpio audio
Environment=PYTHONUNBUFFERED=1

# ── Audio (TTS) ───────────────────────────────────────────────────────────
# A systemd service has NO user login session, so TTS (pyttsx3 → espeak →
# PulseAudio) cannot reach the user's audio by default and stays silent — even
# though it works when you run alas_main manually in a terminal. Point the
# service at the user's PulseAudio socket so it uses the SAME path that works
# manually. Replace 1000 with `id -u <User>` if your uid differs, and run
# `loginctl enable-linger <User>` once so /run/user/1000 exists without a login.
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=PULSE_SERVER=unix:/run/user/1000/pulse/native

# Keep the watcher alive across crashes; it is the only way to start ALAS.
Restart=always
RestartSec=3

# On `systemctl stop`, SIGTERM tells the watcher to also stop any running
# alas_main process group before exiting. Give it time for an ordered shutdown.
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target

# ── Install / manage ──────────────────────────────────────────────────────
#   sudo cp deploy/alas-launcher.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now alas-launcher.service   # start now + at boot
#   journalctl -u alas-launcher -f                      # follow logs
#   journalctl -u alas-launcher -e                      # last logs
#   journalctl -u alas-launcher -b                      # current boot's logs
#   sudo systemctl stop alas-launcher.service           # stop watcher (+ ALAS)
#
# NOTE: the launcher only TOGGLES alas_main. The Jetson never powers off here;
#       cut wall power to shut the device down.

# 1. Kullanıcının uid'sini doğrula (büyük ihtimalle 1000)
id -u alas

# 2. Kullanıcı oturumunu (pulse) login olmadan da ayakta tut
loginctl enable-linger alas

# 3. Güncel servis dosyasını kur (veya /etc/systemd/system/...'i elle düzenle)
sudo cp ~/ALAS_PROJECT/AI-Glasses/deploy/alas-launcher.service /etc/systemd/system/
#    uid 1000 DEĞİLSE iki Environment= satırındaki numarayı düzelt:
#    sudoedit /etc/systemd/system/alas-launcher.service

# 4. Yeniden yükle + başlat
sudo systemctl daemon-reload
sudo systemctl restart alas-launcher

# 5. Butona bas → birkaç saniye içinde "ALAS sistemi hazırlanıyor" duymalısın