# scripts/alas-perf.service
# systemd unit: pin Jetson clocks + fan at boot, before ALAS starts.
# Install:
#   sudo cp scripts/alas-perf.service /etc/systemd/system/
#   sudo systemctl daemon-reload && sudo systemctl enable --now alas-perf
# NOTE: adjust the path below if the repo lives somewhere else on the Jetson.
[Unit]
Description=ALAS Jetson performance setup (nvpmodel/jetson_clocks/fan)
# nvpmodel/fan sysfs nodes are not ready in early boot — order after the
# NVIDIA services and still retry once via the pre-sleep, otherwise the unit
# shows [FAILED] on the boot splash even though a manual run works fine.
After=multi-user.target nvpmodel.service nvfancontrol.service

[Service]
Type=oneshot
ExecStartPre=/bin/sleep 8
ExecStart=/bin/bash /home/alas/ALAS_PROJECT/AI-Glasses/scripts/jetson_setup.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
